# cronjob-mailman-checkdbs.service
[Unit]
Description=mailman checkdbs cronjob

[Service]
ExecStart=/usr/lib/mailman/cron/checkdbs
User=mailman
