# cronjob-mailman-nightly_gzip.service
[Unit]
Description=mailman nightly_gzip cronjob

[Service]
ExecStart=/usr/lib/mailman/cron/nightly_gzip
User=mailman
