# cronjob-mailman-gate_news.service
[Unit]
Description=mailman gate_news cronjob

[Service]
ExecStart=/usr/lib/mailman/cron/gate_news
User=mailman
