# cronjob-mailman-disabled.service
[Unit]
Description=mailman disabled cronjob

[Service]
ExecStart=/usr/lib/mailman/cron/disabled
User=mailman
