# cronjob-mailman-senddigests.service
[Unit]
Description=mailman senddigests cronjob

[Service]
ExecStart=/usr/lib/mailman/cron/senddigests
User=mailman
