# cronjob-mailman-mailpasswds.service
[Unit]
Description=mailman mailpasswds cronjob

[Service]
ExecStart=/usr/lib/mailman/cron/mailpasswds
User=mailman
